# packages/gitlab/gitlab-unicorn.service
[Unit]
Description=GitLab Unicorn Server
Requires=redis.service
After=network.target mysql.service redis.service

[Service]
User=gitlab
Group=gitlab
WorkingDirectory=/usr/share/gitlab/gitlab
Environment=RAILS_ENV=production
SyslogIdentifier=gitlab-unicorn
PIDFile=/var/run/gitlab/unicorn.pid

ExecStart=/usr/bin/unicorn_rails -c /usr/share/gitlab/gitlab/config/unicorn.rb -E production

[Install]
WantedBy=multi-user.target

